# fitzflix_supervisor.ini
[program:fitzflix]
command=/Users/server/Sites/fitzflix/venv/bin/gunicorn -b localhost:8000 --workers 6 --threads 2 --preload fitzflix:app --name fitzflix
directory=/Users/server/Sites/fitzflix
user=server
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

[program:fitzflix-rqscheduler]
command=/Users/server/Sites/fitzflix/venv/bin/rqscheduler
directory=/Users/server/Sites/fitzflix
user=server
autostart=true
autorestart=true
stopsignal=TERM
numprocs=1
process_name=%(program_name)s_%(process_num)02d

[program:fitzflix-sql]
command=/Users/server/Sites/fitzflix/venv/bin/python /Users/server/Sites/fitzflix/supervisor.py fitzflix-sql
directory=/Users/server/Sites/fitzflix
user=server
autostart=true
autorestart=true
stopsignal=TERM
numprocs=1
process_name=%(program_name)s_%(process_num)02d

[program:fitzflix-user-request]
command=/Users/server/Sites/fitzflix/venv/bin/python /Users/server/Sites/fitzflix/supervisor.py fitzflix-user-request
directory=/Users/server/Sites/fitzflix
user=server
autostart=true
autorestart=true
stopsignal=TERM
numprocs=1
process_name=%(program_name)s_%(process_num)02d

[program:fitzflix-import]
command=/Users/server/Sites/fitzflix/venv/bin/python /Users/server/Sites/fitzflix/supervisor.py fitzflix-import fitzflix-file-operation
directory=/Users/server/Sites/fitzflix
user=server
autostart=true
autorestart=true
stopsignal=TERM
numprocs=4
process_name=%(program_name)s_%(process_num)02d

[program:fitzflix-transcode]
command=/Users/server/Sites/fitzflix/venv/bin/python /Users/server/Sites/fitzflix/supervisor.py fitzflix-transcode fitzflix-import fitzflix-file-operation
directory=/Users/server/Sites/fitzflix
user=server
autostart=true
autorestart=true
stopsignal=TERM
numprocs=1
process_name=%(program_name)s_%(process_num)02d

[program:fitzflix-file-operation]
command=/Users/server/Sites/fitzflix/venv/bin/python /Users/server/Sites/fitzflix/supervisor.py fitzflix-file-operation fitzflix-import
directory=/Users/server/Sites/fitzflix
user=server
autostart=true
autorestart=true
stopsignal=TERM
numprocs=5
process_name=%(program_name)s_%(process_num)02d
